# worker/supervisord.d/laravel-php-fpm.conf.example
[program:laravel-php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
numprocs=1