# web/supervisord.d/laravel-swoole.conf.example
[program:laravel-swoole]
process_name=%(program_name)s_%(process_num)02d
command= php /var/www/artisan swoole:http restart
autostart=true
autorestart=true
numprocs=1