# publish/deploy/shared/web/laravel-nginx.conf.example
[program:nginx]
process_name=%(program_name)s_%(process_num)02d
command=nginx
autostart=true
autorestart=true
numprocs=1