# publish/deploy/shared/worker/laravel-schedule.conf.example
[program:laravel-schedule]
process_name=%(program_name)s_%(process_num)02d
command=sh -c "php /var/www/artisan schedule:run >> /dev/null 2>&1; sleep 60"
autostart=true
autorestart=true
numprocs=1